# features/mesh-agent/system_files/etc/systemd/system/meshagent.service
[Unit]
Description=MeshCentral Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/meshagent
EnvironmentFile=-/etc/meshagent/meshagent.env
ExecStartPre=/usr/libexec/meshagent-bootstrap.sh
ExecStart=/usr/bin/bash -lc 'exec /opt/meshagent/meshagent >>/var/log/meshagent/meshagent.log 2>&1'
Restart=always
RestartSec=10
LogsDirectory=meshagent
LogsDirectoryMode=0755

[Install]
WantedBy=multi-user.target