# features/mesh-agent/system_files/etc/systemd/system/meshagent.service
[Unit]
Description=MeshCentral Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/meshagent
EnvironmentFile=-/etc/meshagent/meshagent.env
ExecStartPre=/usr/libexec/meshagent-bootstrap.sh
ExecStart=/opt/meshagent/meshagent
Restart=always
RestartSec=10
LogsDirectory=meshagent
LogsDirectoryMode=0755
StandardOutput=append:/var/log/meshagent/meshagent.log
StandardError=append:/var/log/meshagent/meshagent.log

[Install]
WantedBy=multi-user.target