# features/reset-users/system_files/etc/systemd/system/user@.service.d/kiosk-reset.conf
[Service]
ExecStopPost=/usr/local/bin/kiosk-reset.sh