# features/idle-monitor/system_files/usr/lib/systemd/user/kiosk-monitor.service
[Unit]
Description=Kiosk Data Protection Monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kiosk-monitor.sh