# features/idle-monitor/system_files/usr/lib/systemd/user/kiosk-monitor.timer
[Unit]
Description=Check kiosk idle status periodically

[Timer]
# Wait 1 minute after user logs in before starting checks
OnStartupSec=1min
# Check every 30 seconds thereafter
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target