# features/reset-users/system_files/etc/systemd/system/reset-on-boot.service
[Unit]
Description=Kiosk recovery reset on boot
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/kiosk-reset.sh

# Run very early, before login screen is used
DefaultDependencies=no
Before=display-manager.service

[Install]
WantedBy=multi-user.target