# features/reset-users/system_files/etc/systemd/system/reset-on-boot.service
[Unit]
Description=Reset Kiosk Home Directory on Boot
DefaultDependencies=no
After=local-fs.target
# CRITICAL: Run before the display manager or auto-login starts
Before=display-manager.service getty@tty1.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/kiosk-reset.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target