# features/idle-monitor/system_files/etc/systemd/user/kiosk-monitor.service
[Unit]
Description=Kiosk Data Protection Monitor

[Service]
Type=simple
ExecStart=/usr/local/bin/kiosk-monitor.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target